# user/.local/share/systemd/user/autorandr.service
[Unit]
Description=Run autorandr at startup

[Service]
Type=oneshot
ExecStart=/usr/bin/autorandr -c

[Install]
WantedBy=i3.target
