# user/.local/share/systemd/user/devilspie2.service
[Unit]
Description=Devil's Pie
After=started.target

[Service]
ExecStart=/usr/bin/devilspie2

[Install]
WantedBy=started.target
