# user/.local/share/systemd/user/i3.target
[Unit]
Description=i3 session has started
BindsTo=graphical-session.target
