# user/.local/share/systemd/user/lock_screen.service
[Unit]
Description=Lock on suspend
Before=sleep.target

[Service]
Type=forking
ExecStart=/bin/bash -c lock_screen

[Install]
WantedBy=sleep.target
