# user/.local/share/systemd/user/picom.service
[Unit]
Description=Launch picom (formerly compton)
PartOf=i3.target

[Service]
ExecStart=/bin/bash -c picom -f

[Install]
WantedBy=i3.target
