# user/.local/share/systemd/user/started.timer
[Timer]
OnActiveSec=10
Unit=started.target
RemainAfterElapse=false
AccuracySec=1s
