# user/.local/share/systemd/user/started.target
[Unit]
Description=User session has started
After=default.target
Requisite=default.target
Conflicts=exit.target
