# user/.local/share/systemd/user/test_resume.service
[Unit]
After=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /home/rdu/difernan/resumed

[Install]
WantedBy=sleep.target
