# user/.local/share/systemd/user/test_suspend.service
[Unit]
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /home/rdu/difernan/suspended

[Install]
WantedBy=sleep.target
