# system/etc/udev/rules.d/81-bt-headset.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/local/bin/a2dp-fix-wrapper"
