# system/etc/udev/rules.d/95-battery.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/battery-combined-udev.sh --update"
